# pibroadguard.service
[Unit]
Description=PiBroadGuard – Broadcast Device Security Assessment
After=network.target

[Service]
Type=simple
User=pibroadguard
WorkingDirectory=/opt/pibroadguard
EnvironmentFile=/opt/pibroadguard/.env
ExecStart=/opt/pibroadguard/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
